# third_party/civetweb/resources/Info.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>CFBundleExecutable</key> <string>Civetweb</string>
  <key>CFBundlePackageType</key> <string>APPL</string>
  <key>CFBundleTypeRole</key> <string>None</string>
  <key>CFBundleIconFile</key> <string>civetweb</string>
  <key>CFBundleIconFiles</key> <array>
    <string>civetweb_16x16.png</string>
    <string>civetweb_22x22.png</string>
    <string>civetweb_32x32.png</string>
    <string>civetweb_64x64.png</string>
  </array>
  <key>LSUIElement</key> <true/>
  <key>RunAtLoad</key> <true/>
  <key>Label</key> <string>com.nofacepress.civetweb</string>
  <key>ProgramArguments</key> <array> </array>
  <key>KeepAlive</key> <true/>
</dict>
</plist>
